# xen/autoscreen-udev.rules
SUBSYSTEM=="xen-backend", KERNEL=="console*", RUN+="/etc/udev/autoscreen.sh"
