# libvirt-qemu/virsh-stats.cron
*/5 * * * * root /usr/libexec/virsh-collect-stats.sh
*/5 * * * * root /usr/libexec/virsh-graph-stats.py
