# libvirt-qemu/virsh-stats.cron
*/5 * * * * root /usr/local/sbin/virsh-collect-stats.sh
