# linux/ddcsync.service
[Unit]
Description=reconfigure monitor inputs for system

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/local/bin/ddcsync.sh
Restart=no
